# debian/freeorion.freeoriond.service
[Unit]
Description=FreeOrion server
After=network.target
OnFailure=unit-status-freeorion@.service

[Service]
User=freeorion
Group=freeorion
Environment="XDG_CONFIG_HOME=/etc" "XDG_DATA_HOME=/var/lib/"
ExecStart=/usr/games/freeoriond --hostless --log-file /var/log/freeorion/freeoriond.log --network.server.ai.max 5 --network.server.human.min 2
RestartSec=10
Restart=on-failure
LimitCORE=infinity

[Install]
WantedBy=multi-user.target
