# debian/freeorion.freeoriond.service
[Unit]
Description=FreeOrion server
After=network.target
OnFailure=unit-status-freeorion@.service

[Service]
User=freeorion
Group=freeorion
ExecStart=/usr/games/freeoriond --hostless --log-file /var/log/freeorion/freeoriond.log --mplobby-max-ai 5 --mplobby-min-human 2
RestartSec=10
Restart=on-failure
LimitCORE=infinity

[Install]
WantedBy=multi-user.target
