# debian/freeorion.freeoriond.socket
[Unit]
Description=FreeOrion server

[Socket] 
ListenStream = 12346 
Accept = no

[Install] 
WantedBy = sockets.target
