# debian/unit-status-freeorion@.service
[Unit]
Description=Unit Status Freeorion Service
After=network.target

[Service]
Type=simple
ExecStart=/bin/true
